# fcitx5/server/pinakey-uinput-server.service
[Unit]
Description=PinaKey uinput backspace server (gõ không gạch chân cho app không hỗ trợ SurroundingText)
Documentation=https://github.com/trananhtung/pinakey
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/pinakey-uinput-server
Restart=on-failure
RestartSec=2
# Tăng cường an toàn: daemon chỉ cần mở /dev/uinput + một socket 0600 trong
# $XDG_RUNTIME_DIR/pinakey (systemd tạo sẵn thư mục 0700 qua RuntimeDirectory, #72).
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
RuntimeDirectory=pinakey
RuntimeDirectoryMode=0700
# #113: KHÔNG để systemd xoá đệ quy thư mục khi unit dừng (mặc định no) — xoá cả
# uinput.lock/uinput.sock đang sống của instance chạy ngoài systemd, làm flock suốt đời
# của PR #107 vô hiệu (flock chỉ loại trừ trên cùng inode) → hai daemon song song.
# logind vẫn dọn toàn bộ /run/user/<uid> khi hết phiên nên không rò rỉ.
RuntimeDirectoryPreserve=yes

[Install]
WantedBy=default.target
